# deploy/vps/kora-backend.service
[Unit]
Description=KORA Backend (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kora/backend
EnvironmentFile=/opt/kora/backend/.env
# --workers 1 : APScheduler tourne en mémoire dans le process FastAPI (pas
# de démon externe). Avec plusieurs workers Uvicorn, chacun démarre sa PROPRE
# instance du planificateur, indépendante des autres — le cycle quotidien de
# 6h se programmait donc deux fois (un par worker), avec un vrai risque de
# double exécution/double publication. Un seul worker élimine le problème à
# la racine ; le trafic de cet outil interne (usage mono-utilisateur) ne
# justifie de toute façon pas plusieurs workers.
ExecStart=/opt/kora/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
